# deploy/systemd/nomads-backup-media.service
[Unit]
Description=Nomads Area media backup

[Service]
Type=oneshot
Environment=PROJECT_DIR=/root/Nomads_Area
Environment=BACKUP_DIR=/var/backups/nomads-area
ExecStart=/root/Nomads_Area/scripts/backup_media.sh
